# deploy/conf/cron.d/nginx-reload-certs
# reload nginx every 6 hours to pick up any updated SSL certificates
9 */6 * * * www-data nginx -s reload
